# systemd/document-manager.service
[Unit]
Description=Dokumentenverwaltungssystem
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/OrganisationsAI
Environment="PATH=/home/pi/OrganisationsAI/venv/bin"
ExecStart=/home/pi/OrganisationsAI/venv/bin/python3 main.py
Restart=always
RestartSec=10

# Logs
StandardOutput=append:/var/log/document-manager/app.log
StandardError=append:/var/log/document-manager/error.log

[Install]
WantedBy=multi-user.target
